# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx.err.log
priority=10

[program:backend]
command=python run.py
directory=/app/backend
autostart=true
autorestart=true
environment=FLASK_HOST="0.0.0.0",FLASK_PORT="5000"
stdout_logfile=/app/logs/backend.log
stderr_logfile=/app/logs/backend.err.log
priority=20

[program:ai]
command=python server.py
directory=/app/ai_service
autostart=true
autorestart=true
environment=AI_HOST="0.0.0.0",AI_PORT="8000"
stdout_logfile=/app/logs/ai.log
stderr_logfile=/app/logs/ai.err.log
priority=30
